# etc/systemd/system/ejabberd.service.d/50-OutputLogging.conf
[Service]
StandardOutput=append:/var/log/ejabberd/stdout.log
StandardError=append:/var/log/ejabberd/stderr.log
