# etc/rc.local
#!/bin/sh

for d in ngx-fcgi-cache bg3wiki-cache
do
	dir=/run/shm/$d
	mkdir -p "$dir"
	chmod 700 "$dir"
	chown www-data:www-data "$dir"
done
